# amdgpu-fancontrol.service
[Unit]
Description=AMD GPU Fan Control

[Service]
Type=simple
ExecStart=/usr/local/sbin/amdgpu-fancontrol

[Install]
WantedBy=multi-user.target
